# Prometheus/NodeExporterLaunchDaemon/payload/Library/LaunchDaemons/io.prometheus.NodeExporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.prometheus.NodeExporter</string>
	<key>KeepAlive</key>
	<true/>
	<key>UserName</key>
	<string>_prometheus</string>
	<key>GroupName</key>
	<string>_prometheus</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/NodeExporter/node_exporter</string>
		<string>--collector.boottime</string>
		<string>--collector.cpu</string>
		<string>--collector.diskstats</string>
		<string>--collector.filesystem</string>
		<string>--collector.loadavg</string>
		<string>--collector.meminfo</string>
		<string>--collector.netdev</string>
		<string>--collector.textfile</string>
		<string>--collector.time</string>
		<string>--collector.textfile.directory</string>
		<string>/private/var/spool/node-exporter-textfiles</string>
		<string>--web.listen-address=:9111</string>
	</array>
</dict>
</plist>
